# usbdac_monitor.service
[Unit]
Description=USB DAC Monitor

[Service]
ExecStart=/home/mnmt/pi-fi/usbdac_monitor.sh

[Install]
WantedBy=multi-user.target
